# 99-usb-tethering.rules
ACTION=="bind", SUBSYSTEM=="usb", DRIVER=="rndis_host", RUN+="/usr/local/bin/start-usb-tethering.sh"
